# config/systemd/flick.service
[Unit]
Description=Flick Mobile Shell
Documentation=https://github.com/user/flick

# Start after Phosh is ready (we'll replace it)
After=phosh.service
After=lxc@android.service
After=dbus.socket

# We replace Phosh
Conflicts=phosh.service

[Service]
Type=simple
# Run as root - compositor needs root for hwcomposer, drops privileges for spawned apps
User=root
Group=root

# Use our own runtime dir (not /run/user/UID which gets cleaned up by logind)
RuntimeDirectory=flick
RuntimeDirectoryMode=0755
Environment=XDG_RUNTIME_DIR=/run/flick
Environment=XDG_SESSION_TYPE=wayland
Environment=XDG_CURRENT_DESKTOP=Flick
Environment=EGL_PLATFORM=hwcomposer
Environment=HOME=/home/droidian

# Wait longer for Android container to be fully ready (60 seconds)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do if [ "$(getprop sys.boot_completed)" = "1" ]; then exit 0; fi; sleep 1; done; echo "Android container not ready, continuing anyway"'

# Kill any lingering Wayland clients from previous session
ExecStartPre=-/usr/bin/pkill -9 qmlscene
ExecStartPre=-/usr/bin/pkill -9 Xwayland

# Restart hwcomposer to clear stale frame buffers
ExecStartPre=/bin/sh -c 'ANDROID_SERVICE="(vendor.hwcomposer-.*|vendor.qti.hardware.display.composer)" /usr/lib/halium-wrappers/android-service.sh hwcomposer stop || true'
ExecStartPre=/bin/sleep 2

# Ensure state directory exists and has correct permissions
ExecStartPre=/bin/mkdir -p /home/droidian/.local/state/flick
ExecStartPre=/bin/chown -R droidian:droidian /home/droidian/.local/state/flick

# Start hwcomposer with retry loop (try 3 times with 5 second gaps)
ExecStartPre=/bin/sh -c 'for i in 1 2 3; do ANDROID_SERVICE="(vendor.hwcomposer-.*|vendor.qti.hardware.display.composer)" /usr/lib/halium-wrappers/android-service.sh hwcomposer start && exit 0; sleep 5; done; echo "hwcomposer failed to start after retries"'
ExecStartPre=/bin/sleep 2

# Set default audio sink (ensures audio routing works after boot)
ExecStartPre=/bin/sh -c 'sudo -u droidian XDG_RUNTIME_DIR=/run/user/32011 pactl set-default-sink sink.primary_output || true'

# Main compositor
ExecStart=/home/droidian/Flick/shell/target/release/flick

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flick

# Restart on failure
Restart=on-failure
RestartSec=5

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
