# config/systemd/flick-phone-helper.service
[Unit]
Description=Flick Phone Helper Daemon
Documentation=https://github.com/user/flick
After=ofono.service
BindsTo=flick.service
After=flick.service

[Service]
Type=simple
# Needs root for oFono D-Bus access
User=root

# Clean up stale files before starting
ExecStartPre=/bin/rm -f /tmp/flick_phone_status /tmp/flick_phone_cmd

# Force 2G mode for voice calls (VoLTE not supported on most Droidian devices)
ExecStartPre=-/bin/sh -c 'sleep 5 && mmcli -m 0 --set-allowed-modes="2g" || true'

ExecStart=/usr/bin/python3 /home/droidian/Flick/apps/phone/phone_helper.py daemon

StandardOutput=journal
StandardError=journal
SyslogIdentifier=flick-phone

Restart=on-failure
RestartSec=5

[Install]
WantedBy=flick.service
